# deploy/launchd/com.hushline.social.linkedin.verified-user.weekly.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.hushline.social.linkedin.verified-user.weekly</string>

  <key>ProgramArguments</key>
  <array>
    <string>__REPO_DIR__/scripts/run_verified_user_weekly_linkedin_launchd.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__REPO_DIR__</string>

  <key>UserName</key>
  <string>__USER_NAME__</string>

  <key>GroupName</key>
  <string>__GROUP_NAME__</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>__HOME_DIR__</string>
    <key>HUSHLINE_SOCIAL_ENV_FILE</key>
    <string>__ENV_FILE__</string>
    <key>HUSHLINE_SOCIAL_LAUNCHD_SCOPE</key>
    <string>daemon</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>ProcessType</key>
  <string>Background</string>

  <key>RunAtLoad</key>
  <false/>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>1</integer>
    <key>Hour</key>
    <integer>12</integer>
    <key>Minute</key>
    <integer>10</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>__REPO_DIR__/logs/verified-user-weekly-linkedin.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>__REPO_DIR__/logs/verified-user-weekly-linkedin.stderr.log</string>
</dict>
</plist>
