# deploy/launchd/com.hushline.social.linkedin.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.hushline.social.linkedin.daily</string>

  <key>ProgramArguments</key>
  <array>
    <string>__REPO_DIR__/scripts/run_daily_linkedin_launchd.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__REPO_DIR__</string>

  <key>RunAtLoad</key>
  <false/>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>__REPO_DIR__/logs/linkedin-daily.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>__REPO_DIR__/logs/linkedin-daily.stderr.log</string>
</dict>
</plist>
