# deploy/launchd/com.hushline.social.linkedin.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.hushline.social.linkedin.daily</string>

  <key>ProgramArguments</key>
  <array>
    <string>__REPO_DIR__/scripts/run_daily_linkedin_launchd.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__REPO_DIR__</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>__HOME_DIR__</string>
    <key>HUSHLINE_SOCIAL_ENV_FILE</key>
    <string>__ENV_FILE__</string>
    <key>HUSHLINE_SOCIAL_LAUNCHD_SCOPE</key>
    <string>gui</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Weekday</key>
      <integer>1</integer>
      <key>Hour</key>
      <integer>6</integer>
      <key>Minute</key>
      <integer>10</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>2</integer>
      <key>Hour</key>
      <integer>6</integer>
      <key>Minute</key>
      <integer>10</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>3</integer>
      <key>Hour</key>
      <integer>6</integer>
      <key>Minute</key>
      <integer>10</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>4</integer>
      <key>Hour</key>
      <integer>6</integer>
      <key>Minute</key>
      <integer>10</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>5</integer>
      <key>Hour</key>
      <integer>6</integer>
      <key>Minute</key>
      <integer>10</integer>
    </dict>
  </array>

  <key>StandardOutPath</key>
  <string>__REPO_DIR__/logs/linkedin-daily.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>__REPO_DIR__/logs/linkedin-daily.stderr.log</string>
</dict>
</plist>
